# docker/mihomo/rootfs/etc/s6-overlay/s6-rc.d/mihomo/run
#!/usr/bin/execlineb -P

/mihomo -d /root/.config/mihomo
