# docker/mihomo/rootfs/etc/s6-overlay/s6-rc.d/router/run
#!/usr/bin/execlineb -P

/redirect_pkg_handler
